# configs/php/supervisor.conf
[program:php]
command=php -S 0.0.0.0:80
directory=/var/www
redirect_stderr=true
startretries=9999
startsecs=10